# Player/Systemd/koala-signage-agent.service
[Unit]
Description=Koala Signage Player
After=graphical.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mario
Group=mario
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/1000
WorkingDirectory=/opt/koala-signage
ExecStart=/opt/koala-signage/koala-signage-player --config /etc/koala-signage/config.json
Restart=always
RestartSec=5

[Install]
WantedBy=graphical.target
